# systemd/plc-adjust.service
[Unit]
Description=PLC EdgeFlow eBPF Monitor and Decoder Adjuster
Documentation=https://github.com/http418imateapot/PLC-EdgeFlow-eBPF
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=plcmon
Group=plcmon

# EnvironmentFile allows site-specific overrides without changing this unit.
EnvironmentFile=-/etc/plc-edgeflow/adjust.env

ExecStart=/usr/bin/python3 /opt/plc-edgeflow/adjust.py \
    --serial      ${PLC_SERIAL:-all} \
    --machine_sn  ${PLC_MACHINE_SN:-1} \
    --interval    ${PLC_INTERVAL:-60} \
    --min_delta   ${PLC_MIN_DELTA:-8192} \
    --max_module  ${PLC_MAX_MODULE:-16} \
    --max_unit    ${PLC_MAX_UNIT:-8}

# Restart always so a kernel panic / OOM kill does not permanently stop monitoring.
Restart=always
RestartSec=5
TimeoutStopSec=15

# Grant eBPF and kprobe capabilities without running as full root.
# CAP_BPF        — load and manage BPF programs and maps
# CAP_SYS_ADMIN  — required by BCC for perf_event access on older kernels
# CAP_PERFMON    — modern kernels (5.8+) prefer this over CAP_SYS_ADMIN for perf
AmbientCapabilities=CAP_BPF CAP_SYS_ADMIN CAP_PERFMON
CapabilityBoundingSet=CAP_BPF CAP_SYS_ADMIN CAP_PERFMON

# Harden the service environment.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/sys/fs/bpf /sys/kernel/debug

# Redirect logs to journald (structured JSON lines are preserved as-is).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plc-adjust

[Install]
WantedBy=multi-user.target
