# systemd/plc-decoder@.service
[Unit]
Description=PLC EdgeFlow Decoder for topic %i
Documentation=https://github.com/http418imateapot/PLC-EdgeFlow-eBPF
After=network.target mosquitto.service
Wants=mosquitto.service
BindsTo=plc-adjust.service

[Service]
Type=simple
User=plcmon
Group=plcmon

EnvironmentFile=-/etc/plc-edgeflow/decoder.env

ExecStart=/usr/bin/python3 /opt/plc-edgeflow/decoder.py \
    --topic  %i \
    --broker ${PLC_BROKER:-localhost} \
    --port   ${PLC_BROKER_PORT:-1883} \
    --processor ${PLC_DECODER_PROCESSOR:-lineprotocol} \
    --sqlite-path ${PLC_DECODER_SQLITE_PATH:-/var/lib/plc-edgeflow/points.db} \
    --dlq-path ${PLC_DECODER_DLQ_PATH:-/var/lib/plc-edgeflow/decoder-dlq.db}

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/plc-edgeflow

StandardOutput=journal
StandardError=journal
SyslogIdentifier=plc-decoder@%i

[Install]
WantedBy=multi-user.target
